# infra/tecweb.socket
[Unit]
Description=tecweb socket

[Socket]
ListenStream=/run/tecweb.sock

[Install]
WantedBy=sockets.target